# PREVIOUS-CODE-STACKS/2025-07-13-Version/systemd/ezygallery.service
# ==============================================================
# EzyGallery Gunicorn Service (systemd unit)
# Keeps the EzyGallery Flask app running via Gunicorn.
# Runs under artlistmachine on port 8080 and restarts on failure.
# ==============================================================
[Unit]
Description=EzyGallery Gunicorn Service
After=network.target

[Service]
User=artlistmachine
Group=artlistmachine
WorkingDirectory=/home/artlistmachine/ezygallery
Environment=PATH=/home/artlistmachine/ezygallery/venv/bin
ExecStart=/home/artlistmachine/ezygallery/venv/bin/gunicorn --workers 3 --bind 127.0.0.1:8080 app:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
